# ansible/roles/mysql/templates/mysql-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Run the MySQL backup nightly

[Timer]
OnCalendar=*-*-* {{ '%02d' | format(mysql_backup_hour) }}:{{ '%02d' | format(mysql_backup_minute) }}:00
# Catch up if the host was down at the scheduled time.
Persistent=true

[Install]
WantedBy=timers.target
