# ansible/roles/mysql/templates/mysql-backup.service.j2
# {{ ansible_managed }}
[Unit]
Description=MySQL (Docker) nightly backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
# Optional (leading -): the file is absent when off-site is disabled.
EnvironmentFile=-{{ mysql_base_dir }}/mysql-backup.env
ExecStart=/usr/local/bin/mysql-backup.sh
