# ansible/roles/mysql/templates/mysql-backup.service.j2
# {{ ansible_managed }}
[Unit]
Description=MySQL (Docker) nightly backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mysql-backup.sh
