# ansible/roles/mysql/templates/mysql-backup.service.j2
# {{ ansible_managed }}
[Unit]
Description=MySQL (Docker) nightly backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
# Optional (leading -): off-site is opt-in; the file is absent when disabled.
# Carries the S3 static key + bucket/endpoint/age-recipient for the upload.
EnvironmentFile=-{{ mysql_base_dir }}/mysql-backup.env
ExecStart=/usr/local/bin/mysql-backup.sh
